# openclaw.service
[Unit]
Description=OpenClaw AI Gateway
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/app


ExecStartPre=/usr/local/bin/openclaw-before.sh
ExecStart=openclaw gateway --allow-unconfigured --bind lan

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
#ProtectSystem=strict      # 使 /usr、/boot、/etc 变为只读
#ProtectHome=read-only     # 使 /home、/root、/run/user 变为只读
#ReadWritePaths=/root/openclaw/.openclaw
#ReadWritePaths=/root/openclaw/.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target