# openclaw-in-docker.service
[Unit]
Description=OpenClaw AI Gateway
After=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/app

# Environment variables
Environment="PNPM_HOME=/home/openclaw/.local/share/pnpm"
Environment="PATH=/home/openclaw/.local/bin:/home/openclaw/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/openclaw"
#Environment="XDG_RUNTIME_DIR=/run/user/{{ openclaw_uid_value }}"

# DBus session bus
#Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/{{ openclaw_uid_value }}/bus"

# Start command
ExecStart=openclaw gateway

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/openclaw/.openclaw
ReadWritePaths=/home/openclaw/.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target