# openclaw-in-docker.service
[Unit]
Description=OpenClaw AI Gateway
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/app

# Start command
ExecStart=openclaw gateway --allow-unconfigured

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
#ReadWritePaths=/root/openclaw/.openclaw
#ReadWritePaths=/root/openclaw/.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target